# scripts/volume-analysis.service
[Unit]
Description=Volume Analysis Service - Suspicious Volume Detection
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/media/mra/w/w/amiram/projects/market-ta-generator
ExecStart=/media/mra/w/w/amiram/projects/market-ta-generator/venv/bin/python /media/mra/w/w/amiram/projects/market-ta-generator/volume_analysis_cron.py
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/media/mra/w/w/amiram/projects/market-ta-generator

[Install]
WantedBy=multi-user.target
